# systemd/storycraft.service
[Unit]
Description=StoryCraft AI Service
After=network-online.target

[Service]
User=root
WorkingDirectory=/root/StoryCraft-AI
Environment=PATH=/root/StoryCraft-AI/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=LANG=ru_RU.UTF-8
Environment=LC_ALL=ru_RU.UTF-8
ExecStart=/root/StoryCraft-AI/venv/bin/python main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
